# roles/domain/files/etc/cron.d/domain_mgr_cron
# shell variable for cron
SHELL=/bin/bash
# PATH variable for cron
PATH=/usr/local/bin:/usr/local/sbin:/sbin:/usr/sbin:/bin:/usr/bin
#M   S     T M W  User  Befehl
#-----------------------------------------------------------------
7    11     * * *  root  /bin/bash -c 'cd /root/domain_mgr/ ; /bin/bash prepare.bash'
7    16     * * *  root  /bin/bash -c 'cd /root/domain_mgr/ ; /bin/bash rollover.bash'
7    21     * * *  root  /bin/bash -c 'cd /root/domain_mgr/ ; /bin/bash cleanup.bash'

